# _drafts/Linux/Supervisor/supervisord.d/geth.ini
[program:geth]
directory=/home/wallet/geth-linux-amd64-1.9.15/
user=wallet
command=/home/wallet/geth-linux-amd64-1.9.15/start-geth.sh
stdout_logfile=/tmp/geth-supervisor-out.log
stderr_logfile=/tmp/geth-supervisor-err.log
autostart=true
autorestart=true
startsecs=5
priority=1
stopasgroup=true
killasgroup=true
