# _drafts/Linux/Supervisor/supervisord.d/ngrok.ini
[program:ngrok]
directory=/root/ngrok/
user=root
command=/root/ngrok/start-ngrok.sh
stdout_logfile=/tmp/ngrok-supervisor-out.log
stderr_logfile=/tmp/ngrok-supervisor-err.log
autostart=true
autorestart=true
startsecs=5
priority=1
stopasgroup=true
killasgroup=true
